# install/aurora-ums.service
[Unit]
Description=Aurora User Management Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=aurora
Group=aurora
WorkingDirectory=/opt/aurora/ums
ExecStart=/opt/aurora/ums/aurora-ums-service
Restart=always
RestartSec=3
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
